# services/octoshell-restarter.path
[Path]
PathModified=/var/www/octoshell2/restart.txt

[Install]
WantedBy=multi-user.target